# image/rootfs/etc/systemd/system/crowsnest.service
[Unit]
Description=Crowsnest - Webcam streaming for Mainsail
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=pi
EnvironmentFile=/home/pi/printer_data/systemd/crowsnest.env
ExecStart=/home/pi/crowsnest/crowsnest ${CROWSNEST_ARGS}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
